# laravel-demo/docker/supervisord/php-fpm.conf
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/php-fpm.log
